# etc/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true